# admin/shell/supervisor/monitor_sync_game_logs.conf
[program:monitor_sync_game_logs]
command=/usr/bin/php5.6 /home/vagrant/Code/og_sync/admin/public/index.php cli/sync_game_records/sync_service_start
numprocs=1
process_name=%(program_name)s_%(process_num)03d
directory=/home/vagrant/Code/og_sync/admin
autostart=true
autorestart=true
stdout_logfile=/home/vagrant/Code/og_sync/admin/application/logs/monitor_sync_game_logs_stdout.log.supervisor
stdout_logfile_maxbytes=10MB
stderr_logfile=/home/vagrant/Code/og_sync/admin/application/logs/monitor_sync_game_logs_stderr.log.supervisor
stderr_logfile_maxbytes=10MB
