# admin/shell/supervisor/monitor_sync_one_api.conf
[program:monitor_sync_one_api_72]
command=/usr/bin/php5.6 /home/vagrant/Code/og_sync/admin/public/index.php cli/sync_one_api/sync_one_api_start 72
numprocs=1
process_name=%(program_name)s_%(process_num)03d
directory=/home/vagrant/Code/og_sync/admin
autostart=true
autorestart=true
stdout_logfile=/home/vagrant/Code/og_sync/admin/application/logs/monitor_sync_service_stdout.log.supervisor
stdout_logfile_maxbytes=10MB
stderr_logfile=/home/vagrant/Code/og_sync/admin/application/logs/monitor_sync_service_stderr.log.supervisor
stderr_logfile_maxbytes=10MB
